# restore_kbd_backlight.service
[Unit]
Description=Restore Keyboard Backlight
After=suspend.target

[Service]
Type=oneshot
ExecStart=cp /sys/class/leds/asus::kbd_backlight/brightness /var/lib/misc/kbd_brightness
ExecStop=cp /var/lib/misc/kbd_brightness /sys/class/leds/asus::kbd_backlight/brightness

[Install]
WantedBy=suspend.target
